# units/sonatype-nexus@.service.d/publish-8081.conf
[Service]
# Expose the HTTP port (8081)
Environment='DOCKER_PUBLISH_8081=8081:8081'
